# recipes-core/pending-reboot/files/torizon-ab-pending-reboot.path
[Unit]
Description=Torizon A/B update automatic reboot (watch /run/need-reboot)

[Path]
PathExists=/run/need-reboot

[Install]
WantedBy=multi-user.target
